# .config/systemd/user/bn-refresh.timer
[Unit]
Description=Daily bn refresh-all

[Timer]
OnCalendar=daily
# Catch up a run missed while WSL was shut down at the scheduled time (fires on next boot).
Persistent=true
# Spread the wake so several machines don't all fetch at the same instant.
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
